# paths.list
# auto-apt paths list
#
#<file>	<mark>
#
#<mark>:=
#	! - ignore files under this directory
#	* - standard dirs, no package own this directory
#	other - ignore auto installation
#
home/	!
proc/	!
dev/	!
boot/	*
bin/	*
sbin/	*
tmp/	!
lib/modules/	*
lib/security/	*
lib/	*
etc/cron.d/	*
etc/cron.daily/	*
etc/cron.weekly/	*
etc/cron.monthly/	*
etc/logrotate.d/	*
etc/init.d/	*
etc/menu-methods/	*
etc/X11/Xresources/	*
etc/texmf/	*
etc/pam.d/	*
etc/X11/fonts/misc/	*
etc/X11/fonts/	*
etc/X11/	*
etc/emacs/site-start.d	*
etc/emacs/	*
etc/sound/events/	*
etc/sound/	*
etc/CORBA/servers/	*
etc/CORBA/	*
etc/ppp/ip-up.d/	*
etc/ppp/ip-down.d/	*
etc/ppp/	*
etc/news/	*
etc/	*
usr/doc/	*
usr/local/	!
usr/info/	*
usr/man/	man-db
usr/bin/	*
usr/sbin/	*
usr/X11R6/bin/	*
usr/X11R6/include/X11/bitmaps/	*
usr/X11R6/include/X11/pixmaps/	*
usr/X11R6/include/X11/	*
usr/X11R6/include/	*
usr/X11R6/lib/X11/app-defaults/	*
usr/X11R6/lib/X11/fonts/misc/	*
usr/X11R6/lib/X11/fonts/	*
usr/X11R6/lib/X11/	*
usr/X11R6/lib/	*
usr/X11R6/share/locale/	*
usr/X11R6/share/	*
usr/X11R6/man/	man-db
usr/X11R6/	*
usr/i486-linuxlibc1/include/X11/	*
usr/i486-linuxlibc1/include/	*
usr/i486-linuxlibc1/lib/	*
usr/i486-linuxlibc1/bin/	*
usr/i486-linuxlibc1/	*
usr/games/	*
usr/include/	*
usr/lib/mime/packages/	*
usr/lib/menu/	*
usr/lib/debian-test/tests/	*
usr/lib/debian-test/	*
usr/lib/emacsen-common/packages/install/	*
usr/lib/emacsen-common/packages/remove/	*
usr/lib/games/	*
usr/lib/cgi-bin/	*
usr/lib/GNUstep/Libraries/	*
usr/lib/GNUstep/Headers/	*
usr/lib/GNUstep/Tools/	*
usr/lib/GNUstep/	*
usr/lib/libc5-compat/	*
usr/lib/ruby/1.6/	*
usr/lib/ruby/	*
usr/lib/python1.5/	*
usr/lib/php4/	*
usr/lib/php3/	*
usr/lib/perl5/File/ *
usr/lib/perl5/IO/ *
usr/lib/perl5/5.005/ *
usr/lib/perl5/5.6/ *
usr/lib/perl5/	*
usr/lib/pkgconfig/	*
usr/lib/	*
usr/share/doc/	*
usr/share/doc-base/	*
usr/share/man/	man-db
usr/share/info/	*
usr/share/texmf/fonts/vf	*
usr/share/texmf/fonts/	*
usr/share/texmf/tfm/	*
usr/share/texmf/tex/	*
usr/share/texmf/	*
usr/share/doc-base/	*
usr/share/locale/	*
usr/share/gnome/help/	*
usr/share/gnome/apps/	*
usr/share/gnome/	*
usr/share/games/	*
usr/share/pixmaps/	*
usr/share/emacs/site-lisp/	*
usr/share/emacs/	*
usr/share/perl5/File/	*
usr/share/perl5/IO/	*
usr/share/	*
usr/src/kernel-patches/	*
usr/src/	*
usr/info/	*
usr/man/	man-db
usr/etc/	*
usr/	*
var/spool/mail/	!
var/spool/news/	!
var/tmp/	!
var/lib/games/	*
var/lib/	*
var/log/	*
var/cache/	*
var/games/	*
var/spool/	*
var/state/	*
var/run/	*
var/www/	*
var/	*
